# stt.service
[Unit]
Description=Servidor Vosk Speech-to-Text
After=network.target sound.target

[Service]
User=pi
WorkingDirectory=/home/pi/topibot
ExecStart=/home/pi/topibot/venv/bin/python3 /home/pi/topibot/stt_server.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
